# crontab
*/1 * * * * python3 /email_send.py >> /logfile